# deploy/coach.service
[Unit]
Description=Coach Weightlifting App
After=network.target

[Service]
User=jon
Group=jon
WorkingDirectory=/home/jon/coach2/backend
EnvironmentFile=/home/jon/coach2/deploy/.env
ExecStart=/home/jon/coach2/backend/.venv/bin/gunicorn \
    --workers 2 \
    --bind 0.0.0.0:8000 \
    --timeout 30 \
    config.wsgi:application
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
